# src/data/org.buddiesofbudgie.Services.service.in
[Unit]
Description=Budgie Desktop Services
PartOf=graphical-session.target
ConditionEnvironment=WAYLAND_DISPLAY

[Service]
ExecStart=org.buddiesofbudgie.Services
BusName=org.buddiesofbudgie.Services
Restart=on-failure
Slice=session.slice

[Install]
WantedBy=graphical-session.target